# chapter2-5-your-code-in-multiple-servers/ansible/files/webapp.service
[Unit]
Description=webapp

[Service]
ExecStart=/srv/webapp.py

[Install]
WantedBy=multi-user.target
